# common-packager/src/main/resources/packageing/app/service/servicetemplate-restart.service
[Unit]
Description=Restart ${targetServiceName} Server Service 

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl try-restart ${targetServiceName}

[Install]
WantedBy=timers.target